# lib-print/cups/files/org.cups.cupsd.socket
[Unit]
Description=CUPS Scheduler
PartOf=org.cups.cupsd.service

[Socket]
ListenStream=/run/cups/cups.sock

[Install]
WantedBy=sockets.target
