# app-dev/llama-cpp/files/llama-server.service
# /etc/systemd/system/llama-server.service
[Unit]
Description=Llama Server
Documentation=man:llama-server(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=llama
Group=llama
WorkingDirectory=/opt/llm-models
EnvironmentFile=/etc/llama-server/llama-server.env
Environment=PATH=/opt/llm-models/.local/bin:/usr/local/bin:/usr/bin

ExecStart=/usr/bin/llama-server \
  $LLAMA_SERVER_EXTRA_ARGS \
  --log-verbosity 0 \
  --log-disable

ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=3s

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
